# rootfs/etc/s6-overlay/s6-rc.d/apache2/run
#!/command/with-contenv sh
##!/command/execlineb -P
set -e
. /etc/apache2/envvars
exec apache2 -DNO_DETACH
